# .archive/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty -n -o mohamed_abdelnour --noclear %I $TERM
